# rustdesk-backup.service
[Unit]
Description=RustDesk backup (data directory)

[Service]
Type=oneshot
# Backup the data directory (keys + peer database) using rsync
ExecStart=/usr/bin/rsync -a --delete --chmod=Dg+rx,Fg+r %h/data/ /var/backups/rustdesk/
